# docker/s6-rc.d/migrations/up
#!/command/with-contenv sh
exec s6-setuidgid node /app/docker/run-migrations.sh
