# docker/s6-rc.d/app/finish
#!/command/execlineb -S0
foreground { s6-echo "[app] Application stopped with exit code ${1}" }
s6-svscanctl -t /run/service
