# docker/s6-rc.d/app/run
#!/command/with-contenv sh
echo "[app] Starting application..."
cd /app
exec s6-setuidgid node node dist/index.js
